# deploy/proxmox/meshseer.service
[Unit]
Description=Meshseer FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meshseer
Group=meshseer
WorkingDirectory=/opt/meshseer/app
EnvironmentFile=/etc/meshseer/meshseer.env
ExecStart=/opt/meshseer/app/.venv/bin/python -m meshseer.main
Restart=on-failure
RestartSec=5
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
StateDirectory=meshseer

# Keep the service writable only where it needs to persist state.
ReadWritePaths=/var/lib/meshseer

[Install]
WantedBy=multi-user.target
